# Scripts/99-usbftdi.rules
# For PLMS Developments Consult Cable FTDI FT232 & FT245 USB devices with Vendor ID = 0x0403, Product ID = 0xc7d9 
SYSFS{idVendor}==“0403”, SYSFS{idProduct}==“c7d9”, RUN+=”/sbin/modprobe -q ftdi_sio vendor=0x0403 product=0xc7d9”

